# aiblock.service
[Unit]
Description=AIBlock API Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/blyan
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/root/blyan"
Environment="ALLOWED_ORIGINS=http://165.227.221.225:8000,http://localhost:8000"
Environment="DATABASE_URL=postgresql://blyan:password@localhost/blyandb"
Environment="REDIS_URL=redis://localhost:6379"
Environment="ENABLE_STREAMING=true"
Environment="MAIN_SERVER_URL=http://165.227.221.225:8000"
ExecStart=/usr/bin/python3 /root/blyan/api/server.py
ExecStartPost=/bin/sh -c 'sleep 5 && curl -f http://localhost:8000/health || exit 1'
Restart=always
RestartSec=10
TimeoutStartSec=30
StandardOutput=append:/var/log/aiblock.log
StandardError=append:/var/log/aiblock.error.log

# Health check
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target